# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
logfile_maxbytes=10MB
logfile_backups=3

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/tmp/php-fpm.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/php-fpm-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
autostart=true
autorestart=true
startretries=3
priority=10

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/tmp/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/nginx-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
autostart=true
autorestart=true
startretries=3
priority=20

[program:laravel-scheduler]
command=php artisan schedule:work
stdout_logfile=/tmp/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/tmp/scheduler-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
autostart=true
autorestart=true
startretries=3
priority=30

[program:laravel-queue]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
stdout_logfile=/tmp/queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
stderr_logfile=/tmp/queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
autostart=false
autorestart=true
startretries=3
priority=40
